# etc/systemd/system/ollama.service.d/override.conf
# Caps idle RAM: unload model after 30m idle (matches ollama default,
# explicit so it's not silently held on a 32GB box running
# devcontainers + host Roslyn). `ollama ps` / `ollama stop` to
# inspect/free on demand; bump to e.g. 30m if reload latency annoys.
# num_ctx default is 4096 — too small for opencode's system prompt +
# tool schemas + repo context (model truncates, tool-calling degrades).
# 16k balances headroom vs RAM on a 32GB box; qwen2.5-coder takes 32k+.
[Service]
Environment="OLLAMA_KEEP_ALIVE=30m"
Environment="OLLAMA_CONTEXT_LENGTH=16384"
