# etc/systemd/system/ollama.service.d/override.conf
# Caps idle RAM: unload model after 5m idle (matches ollama default,
# explicit so it's not silently held on a 32GB box running
# devcontainers + host Roslyn). `ollama ps` / `ollama stop` to
# inspect/free on demand; bump to e.g. 30m if reload latency annoys.
[Service]
Environment="OLLAMA_KEEP_ALIVE=5m"
